# appliance/utils/bouncer-rules.service
[Unit]
Description=Bouncer Firewall Rules Scheduler Service

[Service]
Type=oneshot
ExecStart=/opt/bouncer/rules.sh
