# appliance/utils/bouncer-rules.service
[Unit]
Description=Bouncer Firewall Rules Scheduler Service

[Service]
Type=oneshot
ExecStart=/home/marcus/firewall/rules.sh
